# damx-installer/DAMX-Daemon.service
[Unit]
Description=DAMX Daemon - Enforce Fan & RGB Settings on Login
After=graphical-session.target

[Service]
ExecStart=%h/Div-Acer-Manager-Max/DAMX.py --daemon
Restart=on-failure
RestartSec=5
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=default.target
